# board/common/rootfs-overlay/etc/systemd/systemd-repart.service.d/overlay.conf
[Unit]
After=dev-mmcblk0.device
Requires=dev-mmcblk0.device

[Service]
ExecStart=/usr/bin/systemd-repart --dry-run=no --discard=no /dev/mmcblk0
