# board/common/rootfs-overlay/etc/systemd/system/rauc.service.d/override.conf
[Unit]
Wants=systemd-time-wait-sync.service
After=systemd-time-wait-sync.service
